# P3M/supervisor/ovpn-supervison.conf
[program:ovpn-supervisor]
command=/usr/sbin/openvpn-python-supervisor
autostart=true
autorestart=true
user=root